# jaydenarch-base/jaydenarch-base.install
post_install() {
	post_upgrade
}
post_upgrade(){
	systemctl enable --now NetworkManager
	echo "LANG=en_US.UTF-8" > /etc/locale.conf
	systemctl enable --now org.cups.cupsd
}
