# jaydenarch-base/jaydenarch-mac.install
post_install(){
	post_upgrade
}
post_upgrade(){
	systemctl enable --now cpupower
	cpupower frequency-set -g powersave
}
